# config/rmngd.service
# Generated reference unit — install via ./scripts/install-rmng.sh (uses rmngd.service.in).
# Defaults assume clone at ~/dev/projects/RMNG-OS. Override with:
#   RMNG_PROJECT_ROOT=/your/path ./scripts/install-rmng.sh

[Unit]
Description=RMNG-OS agent daemon (rmngd) — local execution plane
Documentation=file:%h/dev/projects/RMNG-OS/docs/integrations/operations-usage.md
After=network-online.target
Wants=network-online.target
ConditionPathExists=%h/.cargo/bin/rmngd

[Service]
Type=simple
WorkingDirectory=%h/dev/projects/RMNG-OS
EnvironmentFile=-%h/.rmng/secrets.env
EnvironmentFile=-%h/.config/environment.d/rmng.conf
Environment=RUST_LOG=info
Environment=RMNG_PROJECT_ROOT=%h/dev/projects/RMNG-OS
Environment=PATH=%h/.cargo/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStartPre=%h/.cargo/bin/rmngd --validate
ExecStart=%h/.cargo/bin/rmngd
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=default.target